# rootfs/etc/s6-overlay/s6-rc.d/dbus-system/up
#!/usr/bin/env execlineb

foreground { mkdir -p /var/run/dbus/ }

/usr/bin/dbus-daemon --system
