# rootfs/etc/s6-overlay/s6-rc.d/xauth/up
#!/usr/bin/env sh

export HOME="/home/${USERNAME}"

touch ~/.Xauthority
chown "$USERNAME" ~/.Xauthority

/usr/bin/xauth add :0.0 . $(xxd -l 16 -p /dev/urandom)
