# rootfs/etc/s6-overlay/s6-rc.d/xvfb-setup/up
#!/usr/bin/env execlineb

env DISPLAY=:0.0 xrandr --output screen --mode 640x480
