# rootfs/etc/s6-overlay/s6-rc.d/rtkit/run
#!/usr/bin/env execlineb

s6-envdir /run/s6/container_environment

exec /usr/bin/rtkitctl --start
