# rootfs/etc/s6-overlay/s6-rc.d/xvfb-wait/up
#!/usr/bin/env execlineb

/etc/s6-overlay/s6-rc.d/xvfb-wait/up.sh
