# rootfs/etc/s6-overlay/s6-rc.d/xvfb-wait/up
#!/usr/bin/env sh

until xset -q -display :0.0 &>/dev/null; do
    sleep 5
done
