# rootfs/etc/s6-overlay/s6-rc.d/sway-daemon/run
#!/usr/bin/env execlineb

s6-envdir /run/sway/environment

exec sway
