# rootfs/etc/s6-overlay/s6-rc.d/sway-setup/up
#!/usr/bin/env execlineb

s6-envdir /run/sway/environment

foreground { s6-sleep 1; swaymsg output * adaptive_sync on; }
