# rootfs/etc/s6-overlay/s6-rc.d/labwc-daemon/run
#!/usr/bin/env execlineb

s6-envdir /run/labwc/environment

exec seatd-launch /usr/bin/labwc
