# rootfs/etc/s6-overlay/s6-rc.d/kodi/run
#!/usr/bin/env execlineb

s6-envdir /run/kodi/environment

exec /usr/bin/kodi-wayland --standalone
