# rootfs/etc/s6-overlay/s6-rc.d/kodi/up
#!/usr/bin/env execlineb

s6-envdir /run/kodi/environment

swaymsg exec WAYLAND_DISPLAY="${WAYLAND_DISPLAY}" kodi-wayland
