# rootfs/etc/s6-overlay/s6-rc.d/pipewire/up
#!/usr/bin/env execlineb

s6-envdir /run/pipewire/environment

exec /usr/libexec/pipewire-launcher
