# rootfs/etc/s6-overlay/s6-rc.d/sunshine/run
#!/usr/bin/env execlineb

setsid /root/.nix-profile/bin/sunshine
