# rootfs/etc/s6-overlay/s6-rc.d/ssh-daemon/run
#!/usr/bin/env execlineb

exec /usr/sbin/sshd
