# rootfs/etc/s6-overlay/s6-rc.d/weston-configure/up
#!/usr/bin/env execlineb

foreground { mkdir -p /var/freerdp/certs/ }

winpr-makecert \
    -rdp \
    -path /var/freerdp/certs/ \
    -n server
