# rootfs/etc/s6-overlay/s6-rc.d/ssh-configure/up
#!/usr/bin/env execlineb

foreground { /usr/bin/ssh-keygen -A }
