# rootfs/etc/s6-overlay/s6-rc.d/weston-daemon/run
#!/usr/bin/env execlineb

s6-envdir /run/s6/container_environment

exec env DISPLAY=:0.0 /usr/bin/weston --config=/etc/weston/weston.ini
