# rootfs/etc/s6-overlay/s6-rc.d/weston-daemon/run
#!/usr/bin/env execlineb

s6-envdir /run/weston/environment

exec /usr/bin/weston \
    --backend=wayland \
    --display="$WAYLAND_DISPLAY"
