# rootfs/etc/s6-overlay/s6-rc.d/wayvnc/run
#!/usr/bin/env execlineb

s6-envdir /run/wayvnc/environment

exec wayvnc \
    -C /etc/wayvnc/config \
    --gpu \
    --max-fps=120 \
    0.0.0.0
