# rootfs/etc/s6-overlay/s6-rc.d/novnc/run
#!/usr/bin/env execlineb

exec bash /opt/novnc/utils/novnc_proxy --vnc 127.0.0.1:5900
