# rootfs/etc/s6-overlay/s6-rc.d/dbus-session/up
#!/usr/bin/env sh

eval "export $(/usr/bin/dbus-launch)"

printf "$DBUS_SESSION_BUS_ADDRESS" > /run/s6/container_environment/DBUS_SESSION_BUS_ADDRESS
printf "$DBUS_SESSION_BUS_PID"     > /run/s6/container_environment/DBUS_SESSION_BUS_PID
