# rootfs/etc/s6-overlay/s6-rc.d/xwayland-setup/up
#!/usr/bin/env execlineb

s6-envdir /run/s6/container_environment

xrandr --output "XWAYLAND0" --mode "1920x1080"
