# rootfs/etc/s6-overlay/s6-rc.d/xwayland-setup/up
#!/usr/bin/env execlineb

env DISPLAY=:1.0 xrandr --output XWAYLAND0 --mode 1920x1080
