# rootfs/etc/s6-overlay/s6-rc.d/xvfb-daemon/run
#!/usr/bin/env execlineb

s6-envdir /run/s6/container_environment

exec env -u DISPLAY /usr/bin/Xvfb :0.0 \
    -nocursor \
    -noreset \
    +extension COMPOSITE \
    +extension DAMAGE \
    +extension DOUBLE-BUFFER \
    +extension GLX \
    +extension RANDR \
    +extension X-Resource \
    +extension XFIXES \
    +extension XVideo \
    +extension XVideo-MotionCompensation
