# rootfs/etc/s6-overlay/s6-rc.d/pulseaudio/run
#!/usr/bin/env execlineb

s6-envdir /run/s6/container_environment

exec /usr/bin/pulseaudio --system --disallow-exit --disable-shm --exit-idle-time=-1
