# rootfs/etc/s6-overlay/s6-rc.d/seatd/run
#!/usr/bin/env execlineb

/usr/bin/seatd
