# rootfs/etc/s6-overlay/s6-rc.d/xrdp-daemon/run
#!/usr/bin/env execlineb

s6-envdir /run/s6/container_environment env -u WAYLAND_DISPLAY /usr/sbin/xrdp --nodaemon
